# packages/genesis39/genesis39/files/uci-defaults/zGenesis39-00003-network-lan
#!/bin/sh
# Copyright (C) 2020 Gregory L. Dietsche
# This program is free software: you can redistribute it and/or modify it under
# the terms of the GNU General Public License version 2 as published by the
# Free Software Foundation. This program is distributed in the hope that it
# will be useful, but WITHOUT ANY WARRANTY; without even the implied warranty
# of MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.
# You should have received a copy of the GNU General Public License along with
# this program. If not, see http://www.gnu.org/licenses/.
#
# This script should only run once. Without this, check, this script
# would run again after an upgrade
#
. /lib/script-guard.sh
script_guard genesis39 00003_network_lan

. /lib/genesis39-functions.sh
. /lib/functions.sh

# Capture and redirect ALL LAN DNS requests so that this router handles them.
# This makes it very hard for someone to use a DNS server other than the one
# we have configured. For example, trying to use 8.8.8.8 for DNS will no longer
# answered by 8.8.8.8. Instead, this router will reply.
#
genesis39_force_zone_dns_to_router lan

exit 0
