# packages/genesis39/genesis39/files/uci-defaults/zGenesis-00010-default-filter
#!/bin/sh
# Copyright (C) 2018 Gregory L. Dietsche
# This program is free software: you can redistribute it and/or modify it under
# the terms of the GNU General Public License version 2 as published by the
# Free Software Foundation. This program is distributed in the hope that it
# will be useful, but WITHOUT ANY WARRANTY; without even the implied warranty
# of MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.
# You should have received a copy of the GNU General Public License along with
# this program. If not, see http://www.gnu.org/licenses/.
#
# This script should only run once. Without this, check, this script
# would run again after an upgrade
#
. /lib/genesis39-functions.sh
run_once_guard genesis39 00010_default_filter

#
# We do NOT want to use the DNS servers advertised by the ISP because
# we want to use the OpenDNS Family Shield DNS servers instead.
#
uci -q batch <<-EOT
  set network.wan.peerdns='0'
  set network.wan6.peerdns='0'

  #DO NOT Allow upstream responses in the 127.0.0.0/8 range, e.g. for RBL services
  delete dhcp.@dnsmasq[0].rebind_localhost
EOT

#
# If dnscrypt-proxy is installed, use it.
# Otherwise, directly use Cisco's Family Shield DNS servers.
#
if command -v dnscrypt-proxy; then
  uci -q batch <<-EOT
    delete dhcp.@dnsmasq[0].server
    add_list dhcp.@dnsmasq[0].server='127.0.0.1#5353'
    add_list dhcp.@dnsmasq[0].server='/pool.ntp.org/8.8.8.8'

    set dnscrypt-proxy.ns1.resolver='cisco-familyshield'
EOT
else
  uci -q batch <<-EOT
    delete dhcp.@dnsmasq[0].server
    add_list dhcp.@dnsmasq[0].server=208.67.222.123
    add_list dhcp.@dnsmasq[0].server=208.67.220.123
EOT
fi

#
# Capture and redirect ALL DNS requests so that this router handles them.
# This makes it very hard for someone to use a DNS server other than the one
# we have configured. For example, trying to use 8.8.8.8 for DNS will no longer
# answered by 8.8.8.8. Instead, this router will reply.
#
genesis39_force_zone_dns_to_router lan

exit 0
