# packages/genesis39/genesis39-opkg/files/uci-defaults/genesis39-opkg
#!/bin/sh
#
# Copyright (c) 2018 Gregory L. Dietsche <Gregory.Dietsche@cuw.edu>
# This is free software, licensed under the MIT License
#
echo src/gz genesis39 http://dl.genesis39.org/packages>> /etc/opkg/customfeeds.conf
exit 0
