# packages/genesis39/genesis39/files/init.d/genesis39-firstboot
#!/bin/sh /etc/rc.common
# Copyright (C) 2018 Gregory L. Dietsche
# This program is free software: you can redistribute it and/or modify it under
# the terms of the GNU General Public License version 2 as published by the
# Free Software Foundation. This program is distributed in the hope that it
# will be useful, but WITHOUT ANY WARRANTY; without even the implied warranty
# of MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.
# You should have received a copy of the GNU General Public License along with
# this program. If not, see http://www.gnu.org/licenses/.
 
START=99
STOP=99
EXTRA_COMMANDS="explain"
EXTRA_HELP="        explain Explains what this service is for"
 
boot() {
  echo "Updating wifi ssid(s)"
  /genesis39/service/firstboot/postfix-wifi-ssids

  echo "Disabling the genesis39-firstboot service"
  /etc/init.d/genesis39-firstboot disable
}

explain(){
  echo This service only executes once on boot. It automatically disables
  echo itsself and should NOT be executed again. It sets up some things
  echo that cannot be configured using /etc/uci-defaults
}
